# mumble/supervisor.conf
[program:mumble]
command=/usr/sbin/murmurd -ini /etc/mumble-server.ini -fg
user=mumble-server
redirect_stderr=true
stdout_logfile=/var/log/mumble-server.log
