# mediacrush/supervisor.celery.conf
[program:celery]
command=/opt/mediacrush/bin/celery worker -A mediacrush
process_name=celery-worker
priority=2
autorestart=true
user=mediacrush
