# mediacrush/supervisor.redis.conf
[program:redis]
command=/usr/sbin/redis-server
process_name=redis
priority=1
user=redis
autorestart=true
