# mediacrush/supervisor.gunicorn.conf
[program:gunicorn]
command=/opt/mediacrush/bin/gunicorn -w 4 app:app
process_name=gunicorn
priority=3
autorestart=true
user=mediacrush
